# deploy/headunit.service
[Unit]
Description=RX7 Headunit Application
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=headunit
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Use installed binary if exists, otherwise fall back to local build
ExecStart=/bin/bash -c 'test -x /usr/bin/headunit && exec /usr/bin/headunit || exec /home/headunit/rx7-headunit/src-tauri/target/release/headunit || exec /home/headunit/rx7-headunit/src-tauri/target/debug/headunit'
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
